# templates/network/sniffing-interface.service.erb
#
# /etc/systemd/system/sniffing-interface-#{interface}.service
#
# Generated by Chef for <%= node[:fqdn] %>
#

[Unit]
Description=Makes interfaces run in promiscuous mode at boot
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/ip link set dev <%= interface %> up
TimeoutStartSec=0
RemainAfterExit=yes

[Install]
WantedBy=default.target
